# power-tweaks/99-50-ac-power.rules
# AC PLUGGED-IN
SUBSYSTEM=="power_supply", ATTR{online}=="1", RUN+="/usr/local/bin/switch-power-settings 0"
