# power-tweaks/99-50-bat-power.rules
# AC NOT PLUGGED-IN
SUBSYSTEM=="power_supply", ATTR{online}=="0", RUN+="/usr/local/bin/switch-power-settings 1"
